# conf/couchdb.conf
#!upstart
description "couchdb server"
author      "greg"

start on startup
stop on shutdown

exec sudo -u couchdb /opt/couchdb-1.0.1/bin/couchdb 2>&1 >> /var/log/couchdb.log
